# files/etc/update-motd.d/20-physcomp
#!/bin/sh
#    10-physcomp - print the localized system identification

printf "\nWelcome to the Carnegie Mellon University Physical Computing Udoo system.\n"

